# .config/systemd/user/rclone-backup.timer
[Unit]
Description=Run Rclone Backup Hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
